# xremap/xremap.service
[Unit]
Description=xremap

[Service]
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=xremap %h/.config/xremap/config.yml

[Install]
WantedBy=default.target
